# mrm_auto_backup.service
[Unit]
Description=MRM RPI Auto Backup

[Service]
ExecStart=/usr/bin/backup.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
